# roles/bashrc/files/etc/bash.bashrc.d/profile.d/99bashrc.sh
# ${bashrcdir}/profile.d/99bashrc.sh
# $Id$

# bash ?
if [ -n "$BASH" ] &&
   [ -e /etc/bash.bashrc ]
then
  . /etc/bash.bashrc
fi

#*eof* 
